# overlay/files-alpine/etc/init.d/c3os-agent
#!/sbin/openrc-run

depend() {
  provide c3os-agent
}

supervisor=supervise-daemon
name="c3os-agent"
command="c3os-agent start"
command_args="--restart"
supervise_daemon_args="--stdout /var/log/c3os/agent.log --stderr /var/log/c3os/agent.log"
pidfile="/run/c3os-agent.pid"
respawn_delay=5
set -o allexport
if [ -f /etc/environment ]; then source /etc/environment; fi
set +o allexport