# overlay/files-alpine/etc/init.d/cos-setup-network
#!/sbin/openrc-run

depend() {
  after net
  provide cos-setup-network
}

start() {
  elemental run-stage network
  eend 0
}