# overlay/files-alpine/etc/init.d/cos-setup-boot
#!/sbin/openrc-run

depend() {
  provide cos-setup-boot
}

start() {
  elemental run-stage boot
  eend 0
}