# overlay/files-alpine/etc/init.d/kairos-agent
#!/sbin/openrc-run

depend() {
  provide kairos-agent
  after cos-setup-network
  use net
}

start() {
    kairos-agent start
    eend $?
}