# files/etc/systemd/system/c3os-setup.service
[Unit]
Description=c3os setup
After=network.target
Wants=network.target
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/c3os setup
[Install]
WantedBy=multi-user.target
