# files-alpine/etc/init.d/c3os-setup
#!/sbin/openrc-run

name=c3os-setup
description="c3os setup"

depend() {
	after cos-setup-network
	use net 
}

start() {
  /usr/bin/c3os setup >> /var/log/c3os-setup.log 2>&1
  eend 0
}