# files-alpine/etc/init.d/c3os-setup
#!/sbin/openrc-run

name=c3os-setup
description="c3os setup"
command="/usr/bin/c3os"
command_args="setup >> /var/log/c3os-setup.log 2>&1"

depend() {
	after  
	use net 
}