# overlay/k8s/rke2/etc/systemd/system/rke2-server.service.d/override.conf
[Service]
ExecStart=
ExecStart=/opt/rke2/bin/rke2 server --disable rke2-ingress-nginx rke2-metrics-server
